# log_analyzer/log_analyzer.conf
[program:cdn_exporter]
directory=/opt/my_website/log_analyzer
command=/opt/my_website_venv/bin/python3 /opt/my_website/log_analyzer/log_analyzer.py
autorestart=true
autostart=true